# bhyve-vms
#!/bin/sh

# PROVIDE: bhyve-vms
# REQUIRE: LOGIN DAEMON network
# KEYWORD: shutdown

. /etc/rc.subr

name="bhyve-vms"
rcvar="${name}_enable"

bhyve_cli_path="/root/script/bhyve-cli/bhyve-cli.sh"
vm_config_base_dir="/usr/local/etc/bhyve-cli/vm.d"

start_bhyve-vms() {
    echo "Initializing bhyve network switches..."
    "${bhyve_cli_path}" switch init

    echo "Starting autostart-enabled bhyve VMs..."
    for VMCONF in "${vm_config_base_dir}"/*/vm.conf; do
        if [ -f "$VMCONF" ]; then
            VMNAME=$(basename "$(dirname "$VMCONF")")
            # Load VM config to get AUTOSTART flag
            # This is a simplified approach, assuming vm.conf only contains key=value pairs
            # and AUTOSTART is always present.
            local AUTOSTART_FLAG=$(grep "^AUTOSTART=" "$VMCONF" | cut -d'=' -f2)
            if [ "$AUTOSTART_FLAG" = "yes" ]; then
                echo "Starting VM: $VMNAME"
                "${bhyve_cli_path}" start "$VMNAME"
            fi
        fi
    done
}

stop_bhyve-vms() {
    echo "Stopping all running bhyve VMs..."
    # Get list of running VMs from bhyve-cli.sh status
    # Assumes 'bhyve-cli.sh status' output format is consistent
    # and VMNAME is the first column, skipping header and separator.
    RUNNING_VMS=$("${bhyve_cli_path}" status | awk 'NR > 2 {print $1}' | grep -v "---")

    if [ -z "$RUNNING_VMS" ]; then
        echo "No bhyve VMs are running."
        return 0
    fi

    for VMNAME in $RUNNING_VMS; do
        echo "Stopping VM: $VMNAME"
        "${bhyve_cli_path}" stop "$VMNAME" --graceful
    done
}

load_rc_config "$name"
run_rc_command "$1"
